# auth-helper/systemd/prusalink-auth.service
[Unit]
Description=PrusaLink Authentication Helper
After=network.target prusalink.service
Requires=network.target

[Service]
Type=simple
User=kkolk
WorkingDirectory=/home/kkolk/prusalink-auth
ExecStart=/home/kkolk/.nvm/versions/node/v24.11.1/bin/node server.js
Restart=always
RestartSec=5
Environment="PRUSALINK_USER=kkolk"
Environment="PRUSALINK_PASS=<set-during-deployment>"
Environment="PORT=3000"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
