# auth-helper/systemd/prusalink-auth.service
[Unit]
Description=PrusaLink Authentication Helper
After=network.target prusalink.service
Requires=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/prusalink-auth
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment="PRUSALINK_USER=maker"
Environment="PRUSALINK_PASS=<will-be-set-during-deployment>"
Environment="PORT=3000"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
